# root/etc/services.d/postfix/run
#!/usr/bin/with-contenv bash

exec s6-setuidgid root /usr/libexec/postfix/master -c /config/postfix -d 2>&1
